# src/keystroke.service
[Unit]
Description=A keystoke count logger daemon

[Service]
Type=forking
ExecStart=/usr/bin/keystroked

[Install]
WantedBy=multi-user.target
